# .docker/.build/supervisord/05-cron.conf
[program:cron]
autorestart=true
command=cron -f
